# test/test-enabled-root/usr/lib/systemd/system/another.service
[Unit]
Description=Another Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Another Service Start"

[Install]
WantedBy=some.target
