# test/test-enabled-root/usr/lib/systemd/system/disabled.service
[Unit]
Description=Disabled Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Disabled Service Start"

[Install]
WantedBy=some.target
