# test/test-enabled-root/usr/lib/systemd/system/maskedstatic.service
[Unit]
Description=Masked Static Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Masked Static Service Start"
