# test/test-enabled-root/usr/lib/systemd/system/maskedruntimestatic.service
[Unit]
Description=Masked Runtime Static Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Masked Runtime Static Service Start"
