# test/test-enabled-root/usr/lib/systemd/system/maskedruntime.service
[Unit]
Description=Masked Runtime Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Masked Runtime Service Start"

[Install]
WantedBy=some.target
