# test/test-enabled-root/usr/lib/systemd/system/masked.service
[Unit]
Description=Masked Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Masked Service Start"

[Install]
WantedBy=some.target
