# test/test-enabled-root/usr/lib/systemd/system/runtime.service
[Unit]
Description=Runtime Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Runtime Service Start"

[Install]
WantedBy=other.target
