# test/test-enabled-root/usr/lib/systemd/system/static.service
[Unit]
Description=Static Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Static Service Start"
