# test/test-enabled-root/usr/lib/systemd/system/templating@three.service
[Unit]
Description=Templating Three Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Templating Three Start"

[Install]
WantedBy=some.target
