# test/test-enabled-root/usr/lib/systemd/system/templating@two.service
[Unit]
Description=Templating Two

[Service]
Type=oneshot
ExecStart=/bin/echo "Templating Two Start"

[Install]
WantedBy=some.target
