# test/test-enabled-root/usr/lib/systemd/system/unique.service
[Unit]
Description=Unique Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Unique Service Start"

[Install]
WantedBy=some.target
